# units/osbuild-composer-proxy.socket
[Socket]
ListenStream=/builds/weldr-api.socket

[Install]
WantedBy=sockets.target
